# profile.d/qyanu-bash-aliases.sh
# some shortcuts
alias lsa="ls -alhi --color=auto"
alias psa="ps -Alf"
alias apgdef="apg -m22 -x22 -a1 -Mncl"
alias grep="grep --color"
alias whatismyip="curl -s https://qyanu.net/whatismyip"
